# supervisord-wiki.conf
[program:wikitten]
command=/root/scripts/start-wiki.sh
numprocs=1
autostart=true
autorestart=false
startsecs=0
priority=200